# belisasari.service
[Unit]
Description=Belisasari Memecoin Hunting Platform
Documentation=https://github.com/your-username/belisasari
After=network.target
Wants=network.target

[Service]
Type=simple
User=belisasari
Group=belisasari
WorkingDirectory=/opt/belisasari
ExecStart=/usr/bin/node start-belisasari-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=belisasari

# Environment variables
Environment=NODE_ENV=production
EnvironmentFile=/opt/belisasari/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/belisasari/logs
ReadWritePaths=/var/log/belisasari

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target
